# docker/conf/supervisord.conf
[supervisord]
nodaemon=true

[program:apache]
priority = 1
command = /usr/local/bin/docker-php-entrypoint apache2-foreground
autorestart=true
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:worker]
priority = 2
stopasgroup = true
user = www-data
directory = /var/www/k-search
command = /usr/local/bin/php bin/console ksearch:data-process:worker
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
